# systemd/videodrop.service
[Unit]
Description=VideoDrop — Telegram/Upload to LXC storage
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/videodrop
ExecStart=/opt/videodrop/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
